# docker/php/supervisor.ini

[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/log/supervisord.pid
nodaemon=true
user = root

[program:php]
command=docker-entrypoint php-fpm -F
user=root

[program:messenger-consume]
command=php /var/www/bin/console messenger:consume async --limit=10
user=root
autostart=true
autorestart=true
numprocs=1
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
